# scripts/crontab
*/5 * * * * cd $HOME/mirage-www-deployment && git pull
